# scripts/systemd/crontab.service
[Unit]
Description=BusyBox crond for SimpleAdmin scheduled tasks
After=multi-user.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /persist/cron
ExecStartPre=/bin/touch /persist/cron/root
ExecStartPre=/bin/chmod 0600 /persist/cron/root
ExecStart=/usr/sbin/crond -f -L /var/log/crontab.log -c /persist/cron
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
