# scripts/init.d/crontab
#!/bin/sh

case "${1:-}" in
    start)
        mkdir -p /persist/cron
        touch /persist/cron/root
        chmod 0600 /persist/cron/root
        /usr/sbin/crond -L /var/log/crontab.log -c /persist/cron
        ;;
    stop)
        killall crond 2>/dev/null || true
        ;;
    restart|force-reload)
        "$0" stop
        sleep 1
        "$0" start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|force-reload}" >&2
        exit 1
        ;;
esac
