# payload/qcport_payload/files/systemd/atcli-perms.service
[Unit]
Description=Set permissions for /dev/smd7 and /dev/smd11
After=dev-smd7.device dev-smd11.device
Wants=dev-smd7.device dev-smd11.device

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'chgrp radio /dev/smd7 /dev/smd11 2>/dev/null || true; chmod 660 /dev/smd7 /dev/smd11 2>/dev/null || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
